# ansible/cms/roles/install/files/cmsLog@.service
[Unit]
Description=cmsLogService
 
[Service]
Environment="SCRIPT_ARGS=-c %i"
ExecStart=/usr/local/bin/cmsLogService $SCRIPT_ARGS
