# ansible/cms/roles/install/files/cmsRanking.service
[Unit]
Description=CMS Ranking system
 
[Service]
ExecStart=/usr/local/bin/cmsRankingWebServer
