# ansible/cms/roles/install/files/cmsResource@.service
[Unit]
Description=cmsResourceService
 
[Service]
Environment="SCRIPT_ARGS=-a %i"
ExecStart=/usr/local/bin/cmsResourceService $SCRIPT_ARGS
